# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=debug
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:dhcpd]
command=/usr/local/bin/kea.sh
autorestart=true